# .config/systemd/user/backup_postgressql.service
[Unit]
Description=PostgreSQL Backup Service
After=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pg_backup.sh
