# .config/systemd/user/sit_reminder.service
[Unit]
Description=reminder for sit positions

[Service]
Type=simple
ExecStart=~/.local/bin/reminders

[Install]
WantedBy=default.target
